# chef/cookbooks/nomad/templates/nomad.service.erb
[Unit]
Description=nomad Agent
Documentation=https://www.nomad.io/
Requires=network-online.target
After=network-online.target

[Service]
User=<%= @service_user %>
Group=<%= @service_group %>
# Use ExecStartPre to create data dir if it doesn't exist
ExecStartPre=/bin/chown -R <%= @service_user %>:<%= @service_group %> <%= @data_dir %>
Restart=on-failure
ExecStart=/usr/local/bin/nomad agent -config=<%= @config_dir %>

[Install]
WantedBy=multi-user.target

